# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:desktop]
command=/opt/openclaw-entrypoint.sh
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/desktop.log
stderr_logfile=/var/log/supervisor/desktop-error.log

[program:hermes-webui]
command=/opt/hermes-webui/venv/bin/python /opt/hermes-webui/server.py
directory=/opt/hermes-webui
environment=HERMES_WEBUI_AGENT_DIR="/opt/hermes",HERMES_WEBUI_PYTHON="/opt/hermes-webui/venv/bin/python",HERMES_WEBUI_HOST="0.0.0.0",HERMES_WEBUI_PORT="8787",HERMES_WEBUI_STATE_DIR="/opt/data/.hermes/webui-mvp",HERMES_WEBUI_DEFAULT_WORKSPACE="/opt/data",HERMES_HOME="/opt/data",HOME="/root",PYTHONUNBUFFERED="1",PATH="/usr/local/bin:/usr/bin:/bin"
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/supervisor/hermes-webui.log
stderr_logfile=/var/log/supervisor/hermes-webui-error.log
